# systemd/groq-tunnel.service
[Unit]
Description=SSH Reverse Tunnel to VPS (BEEBOT Groq proxy)
After=network-online.target groq-proxy.service
Wants=network-online.target
Requires=groq-proxy.service

[Service]
Type=simple
User=hive
# Exposes hive's local port 8990 as port 8990 on the VPS
ExecStart=/usr/bin/ssh \
    -N \
    -o ExitOnForwardFailure=yes \
    -R 8990:localhost:8990 \
    beebot-vps
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=groq-tunnel

[Install]
WantedBy=multi-user.target
